# luci/luci-nginxer/root/etc/init.d/nginxer.init
#!/bin/sh /etc/rc.common

START=93
USE_PROCD=1
 
start_service() { 
	/etc/init.d/nginx enabled
	/etc/init.d/uhttpd disable
}

stop_service() {
	/etc/init.d/nginx disable
	/etc/init.d/uhttpd enabled
}