# luci/luci-nginxer/root/etc/uci-defaults/50_luci-nginxer
#!/bin/sh

/etc/init.d/nginxer enable
/etc/init.d/nginxer start

exit 0
