# luci/luci-nginxer/root/etc/uci-defaults/50_luci-nginxer
#!/bin/sh

/etc/init.d/uhttpd disable
/etc/init.d/nginx enabled
/etc/init.d/nginx start
exit 0
